# setkeycodes.service
[Unit]
Description=Set keycodes for keyboard
[Service]
ExecStart=/etc/systemd/system/setkeycodes.sh
[Install]
WantedBy=multi-user.target
